# cronjobs.cron
21 08 * * 1 python3 app/main.py #run_markets
